# systemd/mcctl-watchdog.service
[Unit]
Description=mcctl self-healing watchdog for the remote Minecraft server
Documentation=https://github.com/lonaivdev-cell/minecraft-remote-monitoring
After=network-online.target

[Service]
ExecStart=/usr/bin/mcctl watchdog run
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
