# systemd/mcctl-backup.service
[Unit]
Description=mcctl world backup + rotation

[Service]
Type=oneshot
ExecStart=/usr/bin/mcctl backup create --notify
