# src/mcctl/units/mcctl-metrics.service
[Unit]
Description=mcctl Prometheus textfile exporter (refresh the .prom from metric history)

[Service]
Type=oneshot
ExecStart=/usr/bin/mcctl metrics export
